# debian-testing-dev/debian-testing-dev/supervisor/sshd.conf
[program:sshd]
command		= /usr/sbin/sshd -D -e
redirect_stderr	= true
stdout_logfile	= /var/log/%(program_name)s.log
