# debian-testing-dev/debian-testing-dev/supervisor/startup.conf
[program:sshd]
command		= /usr/local/bin/startup
redirect_stderr	= true
stdout_logfile	= /var/log/%(program_name)s.log
